# init/lodge.service
[Unit]
Description=lodge

[Service]
ExecStart=/usr/sbin/lodge -p -v -r

[Install]
WantedBy=multi-user.target
